# internal/service/templates/runit/log/run
#!/usr/bin/env sh
# timestamp + message
exec svlogd -tt /var/log/mullvad-daemon