# internal/service/templates/s6/log/run
#!/usr/bin/execlineb -P
s6-log --directory /var/log/mullvad-daemon --timestamps --rotate=7
